# autostart/.config/autostart/chrome.desktop
[Desktop Entry]
Name=Chrome
Comment=Launch google-chrome at startup
Exec=google-chrome
Type=Application
Terminal=false
