# autostart/.config/autostart/pomodoro.desktop
[Desktop Entry]
Name=Gnome Pomodoro
Comment=Start a pomodoro at startup
Exec=gnome-pomodoro --start
Type=Application
Terminal=false
